# switch-DEs.service
[Unit]
Description=restart DE

[Service]
User=root
ExecStart=/bin/bash /usr/bin/switch-DEs.sh
Restart=never

[Install]
WantedBy=multi-user.target

