# waybar-hyprland.service
[Unit]
Description=Waybar Service

[Service]
ExecStart=/bin/waybar
Restart=never

[Install]
WantedBy=multi-user.target

